# roles/neo4j/templates/neo4j-smf.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='neo4j'>
	<service name='neo4j' type='service' version='0'>

		<create_default_instance enabled='false'/>
		<single_instance/>

		<dependency name='network' grouping='require_all' restart_on='error' type='service'>
			<service_fmri value='svc:/milestone/network:default'/>
		</dependency>

		<dependency name='filesystem' grouping='require_all' restart_on='error' type='service'>
			<service_fmri value='svc:/system/filesystem/local'/>
		</dependency>

		<method_context project='{{ neo4j.project }}' working_directory='{{ neo4j.data_path }}'>
			<method_credential user='{{ neo4j.user }}' group='{{ neo4j.group }}'/>
			<method_environment>
				<envvar name='JAVA_HOME' value='/opt/local/java/openjdk8'/>
				<envvar name='PATH' value='/opt/local/sbin:/opt/local/bin:/usr/sbin:/usr/bin'/>
			</method_environment>
		</method_context>

		<exec_method name='start' type='method' exec='{{ neo4j.path }}/bin/neo4j start' timeout_seconds='60'/>
		<exec_method name='stop' type='method' exec=':kill' timeout_seconds='60'/>

		<property_group name='startd' type='framework'>
			<propval name='duration' type='astring' value='contract'/>
			<propval name='ignore_error' type='astring' value='core,signal'/>
		</property_group>

		<stability value='Evolving'/>
		<template>
			<common_name>
				<loctext xml:lang='C'>Neo4j Database Server</loctext>
			</common_name>
		</template>
	</service>
</service_bundle>
